# s6-overlay/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root /app/darkhttpd /out/html --port 4774
