# s6-overlay/etc/services.d/darkhttpd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid ${PUID}:${PGID} darkhttpd /out/html --port 4774
